# timeout-bot.service
[Unit]
Description=Timeout Bot for Discord
After=network.target

[Service]
User=cronrunner
WorkingDirectory=/home/cronrunner/timeout-bot
ExecStart=/bin/bash -c 'source /home/cronrunner/timeout-bot/venv/bin/activate && source /home/cronrunner/credentials/timeout-bot.conf && exec python bot.py --token "$DISCORD_TOKEN" --cooldown-channel "$CHANNEL_ID" --cooldown-role "$COOLDOWN_ROLE_ID"'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
