# timeout-bot.service
[Unit]
Description=Kleurplaat Bot for Discord
After=network.target

[Service]
User=cronrunner
WorkingDirectory=/home/cronrunner/timeout-bot

# Laad environment vars uit bestand
EnvironmentFile=/home/cronrunner/credentials/timeout-bot.conf

# Activeer virtuele omgeving en voer bot uit
ExecStart=/bin/bash -c 'source /home/cronrunner/timeout-bot/venv/bin/activate && exec python bot.py'

# Herstart bot bij crash
Restart=always
RestartSec=5

# Logs naar journalctl
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target